# assets/scripts/reset-guest.service
[Unit]
Description=Reset Guest Account from Template
After=systemd-user-sessions.service
Before=sddm.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/reset-guest.sh
RemainAfterExit=no

[Service]
Type=oneshot
ExecStop=/usr/local/bin/reset-guest.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

[Install]
WantedBy=