# assets/services/guest-login.service
[Unit]
Description=Guest Session Login Setup
After=graphical.target
PartOf=graphical.target

[Service]
Type=oneshot
User=guest
# Mute audio in user context
ExecStart=/usr/bin/wpctl set-mute @DEFAULT_AUDIO_SINK@ 1
# Run iceslab as root (needs NOPASSWD sudoers entry or use separate service)
ExecStartPost=/usr/bin/sudo /opt/iceslab/iceslab -u b

[Install]
WantedBy=graphical.target